# lib/dpkg/info/ajenti.conffiles
/etc/ajenti/config.json
/etc/init.d/ajenti
