# lib/dpkg/info/ajenti.prerm
#!/bin/sh
/etc/init.d/ajenti stop
update-rc.d -f ajenti remove > /dev/null
# Automatically added by dh_pysupport
if which update-python-modules >/dev/null 2>&1; then
	update-python-modules -c  ajenti.public
fi
# End automatically added section


