# lib/dpkg/info/apache2.2-common.list
/etc
/etc/cron.daily
/etc/bash_completion.d
/etc/bash_completion.d/apache2.2-common
/etc/ufw
/etc/ufw/applications.d
/etc/ufw/applications.d/apache2.2-common
/etc/apache2
/etc/apache2/sites-enabled
/etc/apache2/conf.d
/etc/apache2/conf.d/other-vhosts-access-log
/etc/apache2/conf.d/localized-error-pages
/etc/apache2/conf.d/security
/etc/apache2/conf.d/charset
/etc/apache2/mods-enabled
/etc/apache2/mods-available
/etc/apache2/mods-available/ident.load
/etc/apache2/mods-available/authz_default.load
/etc/apache2/mods-available/authn_alias.load
/etc/apache2/mods-available/version.load
/etc/apache2/mods-available/disk_cache.conf
/etc/apache2/mods-available/mem_cache.conf
/etc/apache2/mods-available/authn_default.load
/etc/apache2/mods-available/cern_meta.load
/etc/apache2/mods-available/imagemap.load
/etc/apache2/mods-available/disk_cache.load
/etc/apache2/mods-available/mem_cache.load
/etc/apache2/sites-available
/etc/apache2/sites-available/default-ssl
/etc/apache2/sites-available/default
/etc/init.d
/etc/default
/etc/logrotate.d
/var
/var/cache
/var/cache/apache2
/var/log
/var/log/apache2
/var/www
