# lib/dpkg/info/apt-show-versions.conffiles
/etc/bash_completion.d/apt-show-versions
/etc/cron.daily/apt-show-versions
