# lib/dpkg/info/apt-show-versions.list
/.
/etc
/etc/bash_completion.d
/etc/bash_completion.d/apt-show-versions
/etc/cron.daily
/etc/cron.daily/apt-show-versions
/var
/var/cache
/var/cache/apt-show-versions
/usr
/usr/bin
/usr/bin/apt-show-versions
/usr/share
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/apt-show-versions.1p.gz
/usr/share/man/de
/usr/share/man/de/man1
/usr/share/man/de/man1/apt-show-versions.1.gz
/usr/share/man/es
/usr/share/man/es/man1
/usr/share/man/es/man1/apt-show-versions.1.gz
/usr/share/man/fr
/usr/share/man/fr/man1
/usr/share/man/fr/man1/apt-show-versions.1.gz
/usr/share/man/sv
/usr/share/man/sv/man1
/usr/share/man/sv/man1/apt-show-versions.1.gz
/usr/share/man/vi
/usr/share/man/vi/man1
/usr/share/man/vi/man1/apt-show-versions.1.gz
/usr/share/doc
/usr/share/doc/apt-show-versions
/usr/share/doc/apt-show-versions/copyright
/usr/share/doc/apt-show-versions/changelog.gz
