# lib/dpkg/info/initscripts.list
/.
/etc
/etc/default
/etc/default/devpts
/etc/default/halt
/etc/default/bootlogd
/etc/init.d
/etc/init.d/single
/etc/init.d/stop-bootlogd-single
/etc/init.d/killprocs
/etc/init.d/reboot
/etc/init.d/ondemand
/etc/init.d/rc.local
/etc/init.d/skeleton
/etc/init.d/stop-bootlogd
/etc/init.d/umountroot
/etc/init.d/sendsigs
/etc/init.d/urandom
/etc/init.d/umountnfs.sh
/etc/init.d/umountfs
/etc/init.d/halt
/etc/init.d/bootlogd
/sbin
/sbin/fsck.nfs
/bin
/bin/mountpoint
/var
/var/log
/var/log/fsck
/var/lib
/var/lib/initscripts
/var/lib/urandom
/usr
/usr/share
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/initscripts
/usr/share/initscripts
/usr/share/initscripts/default.rcS
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/fsck.nfs.8.gz
/usr/share/man/man5
/usr/share/man/man5/rcS.5.gz
/usr/share/man/man5/halt.5.gz
/usr/share/man/man1
/usr/share/man/man1/mountpoint.1.gz
/usr/share/doc
/usr/share/doc/initscripts
/usr/share/doc/initscripts/copyright
/usr/share/doc/initscripts/README.Debian
/sys
/run
/lib
/lib/init
/lib/init/splash-functions-base
/lib/init/vars.sh
/usr/share/doc/initscripts/changelog.Debian.gz
/lib/init/rw
