# lib/dpkg/info/php5-fpm.list
/.
/usr
/usr/share
/usr/share/php5
/usr/share/php5/sapi
/usr/share/php5/sapi/fpm
/usr/share/doc
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/php5-fpm.8.gz
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/php5-fpm
/usr/share/bug
/usr/share/bug/php5-fpm
/usr/share/bug/php5-fpm/control
/usr/share/bug/php5-fpm/script
/usr/lib
/usr/lib/php5
/usr/lib/php5/20121212
/usr/lib/php5/php5-fpm-checkconf
/usr/sbin
/usr/sbin/php5-fpm
/etc
/etc/php5
/etc/php5/fpm
/etc/php5/fpm/conf.d
/etc/php5/fpm/php-fpm.conf
/etc/php5/fpm/pool.d
/etc/php5/fpm/pool.d/www.conf
/etc/init
/etc/init/php5-fpm.conf
/etc/logrotate.d
/etc/logrotate.d/php5-fpm
/etc/init.d
/usr/share/doc/php5-fpm
/etc/init.d/php5-fpm
