# lib/dpkg/info/plymouth.list
/.
/etc
/etc/init
/etc/init/plymouth-ready.conf
/etc/init/plymouth.conf
/etc/init/plymouth-log.conf
/etc/init/plymouth-upstart-bridge.conf
/etc/init/plymouth-splash.conf
/etc/init/plymouth-stop.conf
/etc/init.d
/sbin
/sbin/plymouthd
/bin
/bin/plymouth
/bin/plymouth-upstart-bridge
/var
/var/spool
/var/spool/plymouth
/var/lib
/var/lib/plymouth
/usr
/usr/share
/usr/share/initramfs-tools
/usr/share/initramfs-tools/scripts
/usr/share/initramfs-tools/scripts/init-bottom
/usr/share/initramfs-tools/scripts/init-bottom/plymouth
/usr/share/initramfs-tools/scripts/init-top
/usr/share/initramfs-tools/scripts/init-top/plymouth
/usr/share/initramfs-tools/scripts/panic
/usr/share/initramfs-tools/scripts/panic/plymouth
/usr/share/initramfs-tools/conf-hooks.d
/usr/share/initramfs-tools/hooks
/usr/share/initramfs-tools/hooks/plymouth
/usr/share/doc
/usr/share/doc/plymouth
/usr/share/doc/plymouth/copyright
/usr/share/doc/plymouth/README.Debian
/lib
/lib/plymouth
/lib/plymouth/themes
/lib/plymouth/themes/details
/lib/plymouth/themes/details/details.plymouth
/lib/plymouth/renderers
/lib/plymouth/renderers/vga16fb.so
/lib/plymouth/renderers/drm.so
/lib/plymouth/renderers/frame-buffer.so
/lib/plymouth/ubuntu_logo.png
/lib/plymouth/script.so
/lib/plymouth/details.so
/etc/init.d/plymouth-ready
/etc/init.d/plymouth-splash
/etc/init.d/plymouth-log
/etc/init.d/plymouth
/etc/init.d/plymouth-stop
/etc/init.d/plymouth-upstart-bridge
/usr/share/doc/plymouth/changelog.Debian.gz
