# lib/dpkg/info/stunnel4.list
/.
/usr
/usr/bin
/usr/bin/stunnel4
/usr/bin/stunnel3
/usr/lib
/usr/lib/stunnel
/usr/lib/stunnel/libstunnel.so
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/stunnel3.8.gz
/usr/share/man/man8/stunnel4.8.gz
/usr/share/man/fr
/usr/share/man/fr/man8
/usr/share/man/fr/man8/stunnel4.8.gz
/usr/share/man/pl
/usr/share/man/pl/man8
/usr/share/man/pl/man8/stunnel4.8.gz
/usr/share/doc-base
/usr/share/doc-base/stunnel4
/usr/share/doc
/usr/share/doc/stunnel4
/usr/share/doc/stunnel4/README.Debian
/usr/share/doc/stunnel4/examples
/usr/share/doc/stunnel4/examples/stunnel.cnf
/usr/share/doc/stunnel4/examples/ca.pl
/usr/share/doc/stunnel4/examples/ca.html
/usr/share/doc/stunnel4/examples/importCA.html
/usr/share/doc/stunnel4/examples/stunnel.conf-sample
/usr/share/doc/stunnel4/examples/importCA.sh
/usr/share/doc/stunnel4/copyright
/usr/share/doc/stunnel4/TODO
/usr/share/doc/stunnel4/NEWS.Debian.gz
/usr/share/doc/stunnel4/NEWS.gz
/usr/share/doc/stunnel4/stunnel.fr.html
/usr/share/doc/stunnel4/stunnel.html
/usr/share/doc/stunnel4/BUGS
/usr/share/doc/stunnel4/contrib
/usr/share/doc/stunnel4/contrib/StunnelConf-0.1.pl
/usr/share/doc/stunnel4/stunnel.pl.html
/usr/share/doc/stunnel4/README
/usr/share/doc/stunnel4/changelog.Debian.gz
/etc
/etc/default
/etc/default/stunnel4
/etc/logrotate.d
/etc/logrotate.d/stunnel4
/etc/stunnel
/etc/stunnel/README
/etc/ppp
/etc/ppp/ip-down.d
/etc/ppp/ip-down.d/0stunnel4
/etc/ppp/ip-up.d
/etc/ppp/ip-up.d/0stunnel4
/etc/init.d
/etc/init.d/stunnel4
/usr/bin/stunnel
/usr/share/man/man8/stunnel.8.gz
