# lib/dpkg/info/unattended-upgrades.list
/.
/usr
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/unattended-upgrade.8.gz
/usr/share/doc
/usr/share/doc/unattended-upgrades
/usr/share/doc/unattended-upgrades/README
/usr/share/doc/unattended-upgrades/copyright
/usr/share/doc/unattended-upgrades/changelog.gz
/usr/share/pyshared
/usr/share/pyshared/unattended_upgrades-0.1.egg-info
/usr/share/unattended-upgrades
/usr/share/unattended-upgrades/20auto-upgrades-disabled
/usr/share/unattended-upgrades/unattended-upgrade-shutdown
/usr/share/unattended-upgrades/20auto-upgrades
/usr/bin
/usr/bin/unattended-upgrade
/usr/lib
/usr/lib/python2.7
/usr/lib/python2.7/dist-packages
/var
/var/log
/var/log/unattended-upgrades
/etc
/etc/pm
/etc/pm/sleep.d
/etc/pm/sleep.d/10_unattended-upgrades-hibernate
/etc/logrotate.d
/etc/logrotate.d/unattended-upgrades
/etc/init.d
/etc/init.d/unattended-upgrades
/etc/apt
/etc/apt/apt.conf.d
/etc/apt/apt.conf.d/50unattended-upgrades
/usr/bin/unattended-upgrades
/usr/lib/python2.7/dist-packages/unattended_upgrades-0.1.egg-info
